# build/temperature-monitor.service
[Unit]
Description=Temperature Monitor Service
After=network.target

[Service]
WorkingDirectory=/home/$USER/temperature-monitor
ExecStart=/usr/local/bin/node index.js
Restart=on-failure
User=$USER

[Install]
WantedBy=multi-user.target
